# bluetrust.service
[Unit]
Description=Bluetrust

[Service]
ExecStart=/usr/bin/bluetrust
NonBlocking=true
Restart=always
