# bluetrust.service
[Unit]
Description=Bluetrust

[Service]
ExecStart=/usr/bin/bluetrust.py
NonBlocking=true
Restart=always

[Install]
WantedBy=default.target
